# debian/DEBIAN/prerm
systemctl disable caddy
systemctl stop caddy